# roles/custom/gitea_backup/templates/systemd/gitea-backup.timer.j2
[Unit]
Description=Launches gitea-backup.service periodically

[Timer]
Unit=gitea-backup.service
OnCalendar={{ gitea_backup_schedule }}

[Install]
WantedBy=timers.target
