# roles/custom/gitea_backup/templates/systemd/gitea-backup.service.j2
[Unit]
Description=Performs a Gitea backup

[Service]
Type=oneshot
ExecStart={{ gitea_backup_bin_path }}/backup
